# container/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0770
chown=root:agent

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
logfile_backups=3
loglevel=info

[program:mysqld]
command=mysqld --datadir=/var/lib/mysql --user=mysql --bind-address={{MYSQL_BIND}}
priority=100
autorestart=unexpected
exitcodes=0
startsecs=10
startretries=3
stderr_logfile=/var/log/supervisor/mysql.err.log
stdout_logfile=/var/log/supervisor/mysql.out.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
user=root
stopasgroup=true
killasgroup=true
stopsignal=TERM
stopwaitsecs=15

[program:smbd]
command=/usr/sbin/smbd --foreground --no-process-group
autorestart=unexpected
exitcodes=0
stderr_logfile=/var/log/samba/samba.err.log
stdout_logfile=/var/log/samba/samba.out.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
startsecs=5
startretries=5
user=root
priority=200
stopasgroup=true
killasgroup=true
stopsignal=TERM
stopwaitsecs=10

[program:nmbd]
command=/usr/sbin/nmbd --foreground --no-process-group
autorestart=unexpected
exitcodes=0
stderr_logfile=/var/log/samba/nmbd.err.log
stdout_logfile=/var/log/samba/nmbd.out.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
startsecs=5
startretries=5
user=root
priority=200
stopasgroup=true
killasgroup=true
stopsignal=TERM
stopwaitsecs=10

[program:paper]
directory=/minecraft
command=java -Xms{{MC_MIN_MEMORY}} -Xmx{{MC_MAX_MEMORY}} -XX:+UseG1GC -XX:+ParallelRefProcEnabled -XX:MaxGCPauseMillis=200 -jar paper.jar nogui
; Value substituted by entrypoint.sh: "true" or "false"
autostart={{PAPER_AUTOSTART}}
priority=300
autorestart=unexpected
exitcodes=0
startsecs=15
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
user=minecraft
stopsignal=INT
stopwaitsecs=60
stopasgroup=true
killasgroup=true

[program:cron]
command=/usr/sbin/cron -f
autorestart=unexpected
exitcodes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
user=root
stopasgroup=true
killasgroup=true

[program:agent]
directory=/opt/agent
command=node index.js
; Value substituted by entrypoint.sh: "true" or "false"
autostart={{AGENT_AUTOSTART}}
priority=400
autorestart=unexpected
exitcodes=0
startsecs=5
startretries=3
stderr_logfile=/var/log/supervisor/agent.err.log
stdout_logfile=/var/log/supervisor/agent.out.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
user=agent
environment=NODE_ENV="production"
stopasgroup=true
killasgroup=true
stopsignal=TERM
stopwaitsecs=10

[program:web-portal]
directory=/opt/web-portal
command=node build/index.js
; Value substituted by entrypoint.sh: "true" or "false"
autostart={{WEB_PORTAL_AUTOSTART}}
priority=500
autorestart=unexpected
exitcodes=0
startsecs=5
startretries=5
stderr_logfile=/var/log/supervisor/web-portal.err.log
stdout_logfile=/var/log/supervisor/web-portal.out.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
user=agent
environment=NODE_ENV="production",PORT="%(ENV_WEB_PORTAL_PORT)s",WEB_PORTAL_SESSION_SECRET="%(ENV_WEB_PORTAL_SESSION_SECRET)s",WEB_PORTAL_DB_HOST="%(ENV_WEB_PORTAL_DB_HOST)s",WEB_PORTAL_DB_PORT="%(ENV_WEB_PORTAL_DB_PORT)s",WEB_PORTAL_DB_USER="%(ENV_WEB_PORTAL_DB_USER)s",WEB_PORTAL_DB_PASSWORD="%(ENV_WEB_PORTAL_DB_PASSWORD)s",WEB_PORTAL_DB_NAME="%(ENV_WEB_PORTAL_DB_NAME)s",WEB_PORTAL_ADMIN_USER="%(ENV_WEB_PORTAL_ADMIN_USER)s",WEB_PORTAL_ADMIN_PASSWORD="%(ENV_WEB_PORTAL_ADMIN_PASSWORD)s",WEB_PORTAL_RCON_HOST="%(ENV_WEB_PORTAL_RCON_HOST)s",WEB_PORTAL_RCON_PORT="%(ENV_WEB_PORTAL_RCON_PORT)s",WEB_PORTAL_RCON_PASSWORD="%(ENV_WEB_PORTAL_RCON_PASSWORD)s",WEB_PORTAL_GAME_DB_NAME="%(ENV_WEB_PORTAL_GAME_DB_NAME)s",WEB_PORTAL_MC_DIR="%(ENV_WEB_PORTAL_MC_DIR)s",WEB_PORTAL_QUERY_PORT="%(ENV_WEB_PORTAL_QUERY_PORT)s"
stopasgroup=true
killasgroup=true
stopsignal=TERM
stopwaitsecs=10
