# DUITS_MINECRAFT/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
logfile_backups=3
loglevel=info

[program:mysqld]
command=mysqld --datadir=/var/lib/mysql --user=mysql --bind-address={{MYSQL_BIND}}
autorestart=unexpected
exitcodes=0
stderr_logfile=/var/log/supervisor/mysql.err.log
stdout_logfile=/var/log/supervisor/mysql.out.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
user=root
stopasgroup=true
killasgroup=true
stopsignal=TERM
stopwaitsecs=15

[program:smbd]
command=/usr/sbin/smbd --foreground --no-process-group
autorestart=unexpected
exitcodes=0
stderr_logfile=/var/log/samba/samba.err.log
stdout_logfile=/var/log/samba/samba.out.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
startsecs=5
startretries=5
user=root
priority=200
stopasgroup=true
killasgroup=true
stopsignal=TERM
stopwaitsecs=10

[program:nmbd]
command=/usr/sbin/nmbd --foreground --no-process-group
autorestart=unexpected
exitcodes=0
stderr_logfile=/var/log/samba/nmbd.err.log
stdout_logfile=/var/log/samba/nmbd.out.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
startsecs=5
startretries=5
user=root
priority=200
stopasgroup=true
killasgroup=true
stopsignal=TERM
stopwaitsecs=10

[program:paper]
directory=/minecraft
command=java -Xms{{MC_MIN_MEMORY}} -Xmx{{MC_MAX_MEMORY}} -XX:+UseG1GC -XX:+ParallelRefProcEnabled -XX:MaxGCPauseMillis=200 -jar paper.jar --nogui
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
user=minecraft
stopsignal=INT
stopwaitsecs=60
stopasgroup=true
killasgroup=true

[program:cron]
command=/usr/sbin/cron -f
autorestart=unexpected
exitcodes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
user=root
stopasgroup=true
killasgroup=true

[program:agent]
directory=/opt/agent
command=node index.js
; Value substituted by entrypoint.sh: "true" or "false"
autostart={{AGENT_AUTOSTART}}
autorestart=true
stderr_logfile=/var/log/supervisor/agent.err.log
stdout_logfile=/var/log/supervisor/agent.out.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
user=agent
environment=NODE_ENV="production"
stopasgroup=true
killasgroup=true
stopsignal=TERM
stopwaitsecs=10
